# systemd/sirenwatch-pm.service
[Unit]
Description=SirenWatch PM — Projektmanagement
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/sirenwatch/tools/pm
ExecStart=/usr/bin/python3 app.py
Restart=always
RestartSec=5
Environment=PM_PASSWORD=sirenwatch

[Install]
WantedBy=multi-user.target
